# workspace/finance_system/openclaw_schedule.cron.example
CRON_TZ=Asia/Shanghai
0 7 * * * docker exec OPENCLAW_CONTAINER sh -lc 'export OPENCLAW_HOME=/home/node/.openclaw LANG=C.UTF-8 LC_ALL=C.UTF-8 PYTHONIOENCODING=UTF-8 PYTHONUTF8=1; cd /home/node/.openclaw/workspace/finance_system && python3 schedule_reminder.py morning-summary' >> /home/node/.openclaw/workspace/finance_system/schedule_cron.log 2>&1
* * * * * docker exec OPENCLAW_CONTAINER sh -lc 'export OPENCLAW_HOME=/home/node/.openclaw LANG=C.UTF-8 LC_ALL=C.UTF-8 PYTHONIOENCODING=UTF-8 PYTHONUTF8=1; cd /home/node/.openclaw/workspace/finance_system && python3 schedule_reminder.py due-reminders' >> /home/node/.openclaw/workspace/finance_system/schedule_cron.log 2>&1
0 8 * * * docker exec OPENCLAW_CONTAINER sh -lc 'export OPENCLAW_HOME=/home/node/.openclaw LANG=C.UTF-8 LC_ALL=C.UTF-8 PYTHONIOENCODING=UTF-8 PYTHONUTF8=1; cd /home/node/.openclaw/workspace/finance_system && python3 report_bot.py --mode news' >> /home/node/.openclaw/workspace/finance_system/finance_cron.log 2>&1
0 12 * * * docker exec OPENCLAW_CONTAINER sh -lc 'export OPENCLAW_HOME=/home/node/.openclaw LANG=C.UTF-8 LC_ALL=C.UTF-8 PYTHONIOENCODING=UTF-8 PYTHONUTF8=1; cd /home/node/.openclaw/workspace/finance_system && python3 report_bot.py --mode noon' >> /home/node/.openclaw/workspace/finance_system/finance_cron.log 2>&1
